# chef/cookbooks/baseconfig/files/default/rc.local
#!bin/sh
sleep 10
/usr/local/bin/uwsgi --ini /home/ubuntu/project/gamenight/uwsgi.ini --daemonize /var/log/gamenight.log
